# linx-cleanup/linx-cleanup.timer
[Unit]
Description=Run linx-cleanup every hour

[Timer]
OnUnitActiveSec=1h 

[Install]
WantedBy=timers.target
